# network/services/shg-provisioning/files/shg-provisioning
#!/bin/sh /etc/rc.common
#
# Onboarding script for SHG

START=80
STOP=00
USE_PROCD=1

start_service() {
    # FIXME workaround to make device-manager process launched as root by lighttpd
    # TODO see https://redmine.lighttpd.net/projects/1/wiki/docs_modfastcgi (#External Spawning) as a better solution?
    sed -i 's/server.username/#server.username/' /etc/lighttpd/lighttpd.conf
    sed -i 's/server.groupname/#server.groupname/' /etc/lighttpd/lighttpd.conf

    # make uhttpd (LUCI) listen on another port as lighttpd will listen on 80 and 443
    # TODO we should rather change lighttpd configuration
    uci add_list uhttpd.main.listen_http='0.0.0.0:8080'
    uci add_list uhttpd.main.listen_http='[::]:8080'
    uci add_list uhttpd.main.listen_https='0.0.0.0:4443'
    uci add_list uhttpd.main.listen_https='[::]:4443'

    uci commit

    # restart uhttpd (as we changed listening port)
    /etc/init.d/uhttpd restart

    # Reboot
    reboot
}

