# network/services/shg-logger/files/shg-logger.init
#!/bin/sh /etc/rc.common

# Make sure we start after SPIN / MQTT
# NOTE - start MUST be uppercase, and the 
# value MUST be lower than 100
# jb - 2021/06/22
START=94
USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command /usr/bin/python3 /etc/shg/shg-logger.py
    procd_close_instance
}
