# network/services/shg-device-manager/files/shg-device-manager.init
#!/bin/sh /etc/rc.common
#
# Copyright (C) 2020 CIRALabs
# revision 4

START=85 # After mud-supervisor
USE_PROCD=1

start_service() {
    procd_open_instance
    procd_set_param command lua /usr/lib/lua/mud-controller/init.lua
    procd_close_instance
}
